# templates/tomcat.init.erb
#!/bin/sh
# chkconfig: 2345 20 80
# description: Handle tomcat instance <%= name %>

# file managed by puppet

. <%= basedir %>/bin/setenv.sh

# For SELinux we need to use 'runuser' not 'su'
if [ -x "/sbin/runuser" ]; then
    SU="/sbin/runuser"
else
    SU="su"
fi

case "$1" in
  start)
    $SU tomcat -c "umask 0002; $CATALINA_HOME/bin/catalina.sh start"
    RETVAL=$?
    ;;

  stop)
    $SU tomcat -c "$CATALINA_HOME/bin/catalina.sh stop -force"
    RETVAL=$?
  ;;

  restart)
    $SU tomcat -c "$CATALINA_HOME/bin/catalina.sh stop -force"
    sleep 5
    $SU tomcat -c "umask 0002; $CATALINA_HOME/bin/catalina.sh start"
    RETVAL=$?
  ;;

esac

exit $RETVAL
