# templates/tomcat.init.erb
#!/bin/sh
# chkconfig: 2345 20 80
# description: Handle tomcat instance <%= name %>

# file managed by puppet

. <%= basedir %>/bin/setenv.sh

case "$1" in
  start)
    su tomcat -c "umask 0002; $CATALINA_HOME/bin/catalina.sh start"
    RETVAL=$?
    ;;

  stop)
    su tomcat -c "$CATALINA_HOME/bin/catalina.sh stop -force"
    RETVAL=$?
  ;;

  restart)
    su tomcat -c "$CATALINA_HOME/bin/catalina.sh stop -force"
    sleep 5
    su tomcat -c "umask 0002; $CATALINA_HOME/bin/catalina.sh start"
    RETVAL=$?
  ;;

esac

exit $RETVAL
