# templates/tomcat.init.erb
#!/bin/sh
# chkconfig: 2345 20 80
# description: Handle tomcat instance <%= name %>

# file managed by puppet

. /lib/lsb/init-functions

SHUTDOWN_WAIT=30

export JAVA_HOME=<%= javahome %>
export CATALINA_HOME=<%= catalinahome %>
export CATALINA_BASE=<%= basedir %>
export CATALINA_PID=<%= basedir %>/temp/tomcat.pid

# For SELinux we need to use 'runuser' not 'su'
if [ -x "/sbin/runuser" ]; then
    SU="/sbin/runuser"
else
    SU="su"
fi

start() {
  isrunning

  if [ "$?" = 0 ]; then
    echo "tomcat-<%= name %> already running"
    return 0
  fi

  # Source JVM parameters
  . <%= basedir %>/bin/setenv.sh

  $SU tomcat -c "umask 0002; $CATALINA_HOME/bin/catalina.sh start" > /dev/null
}

stop() {
  isrunning

  if [ "$?" = 1 ]; then
    echo "tomcat-<%= name %> already stopped"
    rm -f $CATALINA_PID # remove pidfile if still around
    return 0
  fi

  echo -n "Waiting for tomcat-<%= name %> to exit (${SHUTDOWN_WAIT} sec.): "

  count=0
  until [ "$pid" = "" ] || [ $count -gt $SHUTDOWN_WAIT ]; do
    $SU tomcat -c "$CATALINA_HOME/bin/catalina.sh stop -force" > /dev/null
    findpid

    echo -n "."
    sleep 3
    count=$((count+3))
  done

  echo ""

  if [ "$count" -gt "$SHUTDOWN_WAIT" ]; then
    echo "Forcing tomcat-<%= name %> to stop"
    /bin/kill -9 $pid && sleep 5
  fi

  # check if tomcat is still around, this will be our exit status
  ! isrunning
}

findpid() {
  pid=""
  pid=$(pgrep -U tomcat -f "^$JAVA_HOME/bin/java.*catalina.base=<%= basedir %>")

  # validate output of pgrep
  if ! [ "$pid" = "" ] && ! [ "$pid" -gt 0 ]; then
    log_failure_msg "Unable to determine if tomcat-<%= name %> is running"
    exit 1
  fi
}

isrunning() {
  findpid

  if [ "$pid" = "" ]; then
    return 1
  elif [ "$pid" -gt 0 ]; then
    return 0
  fi
}

case "$1" in
  start)
    start
    RETVAL=$?

    if [ "$RETVAL" = 0 ]; then
      log_success_msg "Started tomcat-<%= name %>"
    else
      log_failure_msg "Not able to start tomcat-<%= name %>"
    fi
    ;;

  stop)
    stop
    RETVAL=$?

    if [ "$RETVAL" = 0 ]; then
      log_success_msg "Stopped tomcat-<%= name %>"
    else
      log_failure_msg "Not able to stop tomcat-<%= name %>"
    fi
  ;;

  restart)
    stop
    sleep 5
    start
    RETVAL=$?

    if [ "$RETVAL" = 0 ]; then
      log_success_msg "Restarted tomcat-<%= name %>"
    else
      log_failure_msg "Not able to restart tomcat-<%= name %>"
    fi
  ;;

  status)
    isrunning
    RETVAL=$?

    if [ "$RETVAL" = 0 ]; then
      echo "tomcat-<%= name %> (pid $pid) is running..."
    else
      echo "tomcat-<%= name %> is stopped"
      RETVAL=3
    fi
  ;;

  *)
    echo "Usage: $0 {start|stop|restart|status}."
  ;;

esac

exit $RETVAL
